# kdeconnect-service/kdeconnect.service
[Unit]
Description=KDE Connect Service for COSMIC
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kdeconnect-service
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
